# teapoy-agent.service
[Unit]
Description=Teapoy Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=onetwopi
WorkingDirectory=/home/onetwopi
EnvironmentFile=/home/onetwopi/.env
ExecStart=/home/onetwopi/agent-teapoy-env/bin/python3 /home/onetwopi/monitor.py --interval 0.5
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
